# AltDaemon/package/DEBIAN/postinst
#!/bin/sh
launchctl load /Library/LaunchDaemons/com.rileytestut.altdaemon.plist
