# AltDaemon/package/Library/LaunchDaemons/com.rileytestut.altdaemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.rileytestut.altdaemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/AltDaemon</string>
	</array>
	<key>UserName</key>
	<string>mobile</string>
	<key>KeepAlive</key>
    <true/>
	<key>RunAtLoad</key>
    <true/>
</dict>
</plist>
